# conf_local/supervisor_celery_flower.ini
[program:celery-flower]
command=/home/boat/.virtualenvs/MxShop/bin/python manage.py celery flower --basic_auth=admin:admin123          ; the program (relative uses PATH, can take args)
directory=/home/boat/django_learn/MxShop                ; directory to cwd to before exec (def no cwd)
environment=PATH="/home/boat/.virtualenvs/MxShop/bin"       ; process environment additions (def no adds)
stdout_logfile=/home/boat/django_learn/MxShop/conf/logs/celery_flower.log
stderr_logfile=/home/boat/django_learn/MxShop/conf/logs/celery_flower_error.log
autostart=true
autorestart=true
stopwaitsecs=60
startretries=3
priority=996

