# conf_local/supervisor_celery_worker.ini
[program:celery-worker]
command=/home/boat/.virtualenvs/MxShop/bin/python manage.py celery worker -l INFO           ; the program (relative uses PATH, can take args)
directory=/home/boat/django_learn/MxShop                ; directory to cwd to before exec (def no cwd)
environment=PATH="/home/boat/.virtualenvs/MxShop/bin"       ; process environment additions (def no adds)
stdout_logfile=/home/boat/django_learn/MxShop/conf/logs/celery_worker.log
stderr_logfile=/home/boat/django_learn/MxShop/conf/logs/celery_worker_error.log
autostart=true
autorestart=true
stopwaitsecs=10
startretries=3
priority=998

